# myservice1.service
[Unit]
Description=My Service

[Service]
Type=simple
ExecStart=/home/aynur/bsm/myscript.sh

[Install]
WantedBy=multi-user.target
